# Docker-Anzo/src/etc/supervisor/supervisor.conf
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid 
logfile = /var/log/supervisor/supervisord.log
childlogdir = /var/log/supervisor

[include]
files = /etc/supervisor/conf.d/*.conf

[supervisorctl]
#serverurl = unix:///dev/shm/supervisor.sock
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor] 
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface  


[unix_http_server]
file = /var/run/supervisor.sock
#file = /dev/shm/supervisor.sock
chmod=0700



[program:anzostudio]
command=/opt/Anzo/Server/AnzoServer start
autostart = true 
autorestart = false
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
startsecs = 0
